# deb-build/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable spdktgt-svr.service
systemctl start spdktgt-svr.service